# install/etc/init.d/pui
#!/sbin/openrc-run

command=/usr/local/bin/pui