# init.d/pui
#!/sbin/openrc-run

command=/usr/bin/pui